# laravel.queue.work.ini
[program:laravel-queue-work]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/laravel-5.7-boilerplate/artisan queue:work --queue[=QUEUE] --tries=3
autostart=true
autorestart=true
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/laravel.queue.work.log
