# laravel.queue.work.ini
[program:tvchannels-queue-work]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/tv-channels-validator/artisan queue:work --queue=tvchannels --tries=3
autostart=true
autorestart=true
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/laravel.queue.work.log
